# templates/simple-program.conf.j2
[program:{{supervisord_program_name}}]

command={{supervisord_program_command}}
autorestart={{supervisord_program_autorestart}}
