# templates/simple-program.conf.j2
[program:{{supervisord_program_name}}]

user={{supervisord_program_username}}
command={{supervisord_program_command}}
autorestart={{supervisord_program_autorestart|default(true)|to_json}}
