# templates/simple-program.conf.j2
[program:{{program-name]

command=/bin/launch-my-daemon
autorestart=true
